# deploy/supervisor/haleelo-worker.conf
; ───────────────────────────────────────────────────────────────────────────
; Supervisor config for the Laravel queue worker (WhatsApp/email jobs, PDFs).
; Place at: /etc/supervisor/conf.d/haleelo-worker.conf
; Then:
;   supervisorctl reread
;   supervisorctl update
;   supervisorctl status haleelo-worker:*
;
; NOTE: With QUEUE_CONNECTION=sync (pilot default) jobs run inline and this
; worker is optional. Switch QUEUE_CONNECTION=database (and run the queue
; migration) when you want async processing, then this worker is required.
; ───────────────────────────────────────────────────────────────────────────

[program:haleelo-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/haleelo/api/artisan queue:work --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/haleelo/worker.log
stopwaitsecs=3600
